# root/etc/uci-defaults/99_luci-theme-desktop
#!/bin/sh
# Set desktop as the default theme
uci -q batch <<-EOF >/dev/null
    set luci.main.mediaurlbase=/luci-static/desktop
    commit luci
EOF
exit 0
